# SuspendFix.service
[Unit]
Description=wakeup-fix
 
[Service]
ExecStart=/bin/bash SuspendFix.sh
 
[Install]
WantedBy=multi-user.target
